# cfg/s6-rc.d/postfix/run
#!/command/execlineb -P
fdmove -c 2 1
/setup/scripts/postfix.sh
